# crontab
# m h  dom mon dow   command
0 20 * * 5 python /code/manage.py email
