# crontab
# m h  dom mon dow   command
0 1 * * * python /code/manage.py email
